# .flocks/flockshub/plugins/skills/Anthropic-Cybersecurity-Skills/hunting-for-suspicious-scheduled-tasks/manifest.json
{
  "schemaVersion": "hub.plugin.v1",
  "id": "hunting-for-suspicious-scheduled-tasks",
  "type": "skill",
  "name": "hunting-for-suspicious-scheduled-tasks",
  "description": "Hunt for adversary persistence and execution via Windows scheduled tasks by analyzing task creation events, suspicious task properties, and unusual execution patterns that indicate T1053.005 abuse.",
  "version": "1.0",
  "author": "mahipal",
  "license": "Apache-2.0",
  "homepage": "https://github.com/mukul975/Anthropic-Cybersecurity-Skills",
  "category": "threat-intel",
  "tags": [
    "windows",
    "ioc"
  ],
  "useCases": [
    "threat-detection"
  ],
  "domains": [
    "security-ops"
  ],
  "capabilities": [
    "llm-agent"
  ],
  "trust": "community",
  "source": {
    "kind": "bundled",
    "path": ".flocks/flockshub/plugins/skills/Anthropic-Cybersecurity-Skills/hunting-for-suspicious-scheduled-tasks"
  },
  "compatibility": {
    "flocks": ">=0.8.0",
    "os": [
      "darwin",
      "linux",
      "windows"
    ]
  },
  "dependencies": {
    "skills": [],
    "tools": [],
    "python": [],
    "external": []
  },
  "permissions": {
    "tools": [],
    "network": false,
    "shell": false,
    "filesystem": "read"
  },
  "risk": {
    "level": "low",
    "reasons": []
  },
  "entrypoints": [
    "SKILL.md"
  ],
  "checksums": {},
  "upstream": {
    "name": "mukul975/Anthropic-Cybersecurity-Skills",
    "url": "https://github.com/mukul975/Anthropic-Cybersecurity-Skills"
  },
  "sourceNotice": "Source: This skill is adapted from the open-source project https://github.com/mukul975/Anthropic-Cybersecurity-Skills."
}
